# uit-tss-toolbox_2.0-1_amd64/etc/cron.d/01-uit-cronjobs
* * * * * root /usr/sbin/uit-sql-refresh