# uit-tss-toolbox_2.0-1_amd64/etc/cron.d/02-uit-tss-cronjobs
* */12 * * * root /usr/sbin/uit-tss-toolbox-backup
* */3 * * * root /usr/sbin/uit-tss-sql-csv
* * * * * root mv /tmp/laptop-report-*.csv /mnt/