# uit-tss-toolbox_2.0-1_amd64/etc/cron.d/02-uit-tss-cronjobs
* */12 * * * root /usr/sbin/uit-tss-toolbox-backup
1 0/3 * * * root /usr/sbin/uit-tss-sql-csv
* * * * * root mv /tmp/laptop-report-job*.csv /mnt/laptop-reports/job-reports
* * * * * root mv /tmp/laptop-report-location*.csv /mnt/laptop-reports/location-reports
* * * * * root mv /tmp/laptop-report-client*.csv /mnt/laptop-reports/client-reports
* * * * * root mv /tmp/laptop-report-server*.csv /mnt/laptop-reports/server-reports