# uit-tss-toolbox_2.0-1_amd64/etc/cron.d/uit-cron
#*/5 * * * * root /usr/sbin/uit-sql-refresh
0 */12 * * * root /usr/sbin/uit-backup
