# uit-tss-toolbox_2.0-1_amd64/etc/cron.d/uit-cron
0 */3 * * * root /usr/sbin/uit-sql-refresh
0 */24 * * * root /usr/sbin/uit-backup
