# uit-tss-toolbox_1.0-1_amd64/etc/cron.d/01-uit-tss-cronjobs
* * * * * root /usr/sbin/uit-tss-sql-refresh