# uit-tss-toolbox_1.0-1_amd64/etc/cron.d/02-uit-tss-cronjobs
*/10 * * * * root /usr/sbin/uit-tss-sql-backup
*/10 * * * * root /usr/sbin/uit-tss-sql-report > /opt/Laptop-Report-$(date --iso).html