# uit-tss-toolbox_1.0-1_amd64/etc/cron.d/02-uit-tss-cronjobs
*/10 * * * * root /usr/sbin/uit-tss-sql-backup
* */3 * * * root /usr/sbin/uit-tss-sql-csv
* * * * * root mv /tmp/laptop-report-*.csv /mnt/