# uit-tss-toolbox_1.0-1_amd64/etc/cron.d/03-uit-tss-cronjobs
@daily root rsync -a --delete /home/hp/ /mnt/
@daily root rsync -a --delete /home/dell/ /mnt/
@daily root rsync -a --delete /home/desktops/ /mnt/