# uit-tss-toolbox_1.0-1_amd64/etc/cron.d/sql-cronjobs
*/10 * * * * root /usr/sbin/uit-tss-sql-backup
* * * * * root /usr/sbin/uit-tss-sql-refresh
*/10 * * * * root /usr/sbin/uit-tss-sql-report > /opt/Laptop-Report-$(date --iso).html
