# uit-tss-toolbox_1.0-1_amd64/etc/cron.d/sql-cronjobs
*/10 * * * * root php /var/lib/uit-tss-sql-backup.php
* * * * * root php /var/lib/uit-tss-sql-refresh.php
